# CA4/puppet/modules/h2/templates/h2.service.erb
[Unit]
Description=H2 Database Server
After=network.target

[Service]
WorkingDirectory=/opt/h2
ExecStart=/usr/bin/java -cp /opt/h2/h2.jar org.h2.tools.Server -tcp -tcpPort 9092 -tcpAllowOthers -baseDir /opt/h2
Restart=always

[Install]
WantedBy=multi-user.target